# deploy/douyin-feishu-watcher.service
[Unit]
Description=Douyin Feishu Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/douyin-feishu-watcher
ExecStart=/opt/douyin-feishu-watcher/.venv/bin/douyin-feishu-watcher serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
